# root/etc/crontabs/abc
# min   hour    day     month   weekday command
*/1     *       *       *       *       sleep $((60 + $RANDOM % 230)); /app/duck.sh ; /app/dynu.sh 2>&1
